# supervisord.conf
[supervisord]
nodaemon=true

[program:memcached]
command=/bin/sh -c "source /etc/memcached/memcached.cfg && memcached -p $port -m $mlimit $allowlru -u $user -P $pidfile"
